# debian/cron.daily
#!/bin/sh

reprepro --confdir /etc/reprepro/debian/ update squeeze
reprepro --confdir /etc/reprepro/debian/ update sid
reprepro --confdir /etc/reprepro/ubuntu/ update oneiric
reprepro --confdir /etc/reprepro/ubuntu/ update precise	