# debian/cron.daily
#!/bin/sh

# update supported debian releases
for d in squeeze sid
do
    reprepro --confdir /etc/reprepro/debian/ -V update $d
done

# update supported ubuntu releases
for u in oneiric precise
do
    reprepro --confdir /etc/reprepro/ubuntu/ -V update $u
done
